# cnproxy/rootfs/etc/s6-overlay/s6-rc.d/nginx-github/finish
#!/command/with-contenv bashio
bashio::log.info "CNProxy nginx-github service stopped."
